# todo-timer.timer
[Unit]
Description=5min timer

[Timer]
OnUnitActiveSec=1200

[Install]
WantedBy=basic.target 
