# todo-timer.timer
[Unit]
Description=5min timer

[Timer]
OnUnitActiveSec=60

[Install]
WantedBy=basic.target 
