# todo-timer.service
[Unit]

[Service]
ExecStart=/usr/bin/reminder.sh

[Install]
WantedBy=default.target
