# todo-timer.service
[Unit]

[Service]
ExecStart=/usr/bin/wall -n "hellllloo" 

[Install]
WantedBy=default.target
